# magento-ce/etc/cron.d/magento
* * * * * root /bin/bash /app/public/cron.sh >> /var/log/cron.log 2>&1
